# hosts/hermes/files/etc/uci-defaults/99-hermes-pppoe-recovery
#!/bin/sh
set -eu

RECOVERY_SCRIPT="/usr/libexec/hermes-pppoe-recovery"
LOG_TAG="hermes-pppoe-recovery"

if [ ! -f "$RECOVERY_SCRIPT" ]; then
  logger -t "$LOG_TAG" "missing recovery script: $RECOVERY_SCRIPT"
  exit 0
fi

/bin/sh "$RECOVERY_SCRIPT" >/tmp/hermes-pppoe-recovery-firstboot.log 2>&1 || true

exit 0
